# systemd_user/.config/systemd/user/wallpaper.service
[Unit]
Description=Set wallpaper for user

[Service]
Type=oneshot
ExecStart=%h/.local/bin/set_wallpaper

[Install]
WantedBy=default.target
